# etc/pm-alert.timer
[Unit]
Description=Fleet alert timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target