# etc/pm-boot-health.service
[Unit]
Description=Boot health scoring
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/pm-boot-health.py

[Install]
WantedBy=multi-user.target