# etc/pm-fleet-agent.service
[Unit]
Description=Fleet telemetry agent

[Service]
Type=oneshot
ExecStart=/opt/pm-fleet-agent.py