# etc/pm-key-rotate.service
[Unit]
Description=Rotate device key

[Service]
Type=oneshot
ExecStart=/opt/pm-rotate-device-key