# etc/storage-sync.service
[Unit]
Description=Periodic Storage Sync

[Service]
Type=oneshot
ExecStart=/bin/sync